# docker/files/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
user = root
autostart = true