# docker/files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/local/sbin/php-fpm
user = root
autostart = true